# packaging/launchd/com.nexusproxy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "https://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.nexusproxy</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/YOUR_USER/.local/bin/nexusproxy</string>
		<string>--config</string>
		<string>/Users/YOUR_USER/.config/nexusproxy/config.json</string>
	</array>
	<key>EnvironmentVariables</key>
	<dict>
		<key>NEXUS_ENV_FILE</key>
		<string>/Users/YOUR_USER/.config/nexusproxy/.env</string>
	</dict>
	<key>WorkingDirectory</key>
	<string>/Users/YOUR_USER/.config/nexusproxy</string>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>StandardOutPath</key>
	<string>/tmp/nexusproxy.out.log</string>
	<key>StandardErrorPath</key>
	<string>/tmp/nexusproxy.err.log</string>
</dict>
</plist>
